# build/lib/systemd/system/galera_watchdog.service
[Unit]
Description=Galera Node Watchdog

[Service]
User=galera_watchdog
Group=galera_watchdog
ExecStart=/opt/bin/galera_watchdog -config /etc/galera_watchdog/galera_watchdog

[Install]
WantedBy=multi-user.target
